# roles/borgbackup/templates/backup.service
[Unit]
Description=Borg Backup Service

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/borgbackup.sh

[Install]
WantedBy=multi-user.target